# mitmproxy/root/etc/services.d/mitmproxy/run
#!/usr/bin/with-contenv bashio
set -e
CONFIG_PATH=/data/options.json
bashio::log.info "Extracting options from ${CONFIG_PATH}..."
OPTIONS=$(bashio::jq "${CONFIG_PATH}" 'with_entries(if .value == "" then empty else . end ) | to_entries[] | ["\(.key)=\(.value)" ] | [.[] | "--set "  + .] |  join(" ") ')
bashio::log.info "Extraction complete."
bashio::log.info "Start mitmproxy with options: $OPTIONS"
HOST=$(bashio::host.hostname)
bashio::log.info "Hostname is ${HOST}"

bashio::log.info "For the onboarding page, configure web client to use proxy ${HOST}:8080, then visit the magic url http://mitm.it which mitmweb will intercept and give you a root certificate download page for your client."
# shellcheck disable=SC2086
mitmweb $OPTIONS \
  --set 'web_host=0.0.0.0' \
  --set 'web_port=9090' \
  --set 'listen_port=8080' \
  --set 'confdir=/data/mitmproxy' \
  --set 'web_open_browser=false' \
  --no-web-open-browser
