# podtube/rootfs/etc/services.d/example/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here

# Declare variables
declare google_api_key

## Get the 'google_api_key' key from the user config options.
google_api_key=$(bashio::config 'google_api_key')

# Since pytube is a cat-and-mouse situation, upgrade on addon start.
bashio::log.info "Checking pytube for new versions..."
pip3 install --upgrade pytube

bashio::log.info "PodTube starting..."
exec python3 /usr/bin/podtube.py $google_api_key 80
